# systemd/asr-transformers.service
[Unit]
Description=ASR Transformers inference service
After=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/code/linux-asr-voicetype
ExecStart=%h/code/linux-asr-voicetype/scripts/start_asr_transformers.sh
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
