# systemd/voicetype.service
[Unit]
Description=VoiceType ASR backend
After=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/code/linux-asr-voicetype
Environment=UV_CACHE_DIR=/tmp/uv-cache
ExecStart=%h/.local/bin/uv run voicetype serve --host 127.0.0.1 --port 8787 --model Qwen/Qwen3-ASR-0.6B
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
